# deploy/systemd/dockscheduler.service
[Unit]
Description=DockScheduler (Flask/gunicorn)
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/soundcheck.tools/DockScheduler
ExecStart=/opt/soundcheck.tools/DockScheduler/.venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5001 \
    app:app
Environment="SECRET_KEY=f764818ae0cc0b7deef1c3c8b63b34e47c13ae0eb4fc865c4e34fc81281e3cfb"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
